# src/usr/lib/systemd/system/lug-sites-sync.timer
[Unit]
Description=LUG@USTC Sites Sync Service Timer

[Install]
WantedBy=multi-user.target

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true
